# bin/saas.service
[Unit]
Description=SleepAsAService

[Service]
Type=simple
ExecStart=/usr/bin/saas loop
